# etc/systemd/system/getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noreset --noclear --nohints --issue-file=/etc/issue.d/issue - ${TERM}
TTYVTDisallocate=no
